# scheduler/onChainSignalling.service
[Unit]
Description = On Chain Signalling Votes API
After = network.target